# supervisord.validate.conf
; Build-time validation config for the WhatsApp bridge Chromium launch.
;
; This is NOT used at runtime. It is invoked once during the Docker build with
; supervisord as the RUN command, so supervisord becomes PID 1 of the build
; container and reproduces the exact runtime condition (PID 1 reaping Chromium's
; child processes) that caused "Failed to launch the browser process: Code: null".
;
; launch-check.js launches Chromium with the production flags, prints
; LAUNCH_CHECK_OK on success, and signals this supervisord (PID 1) to shut down so
; the build step returns promptly. Program output is written to
; /tmp/launch-check.log; a following Dockerfile RUN asserts LAUNCH_CHECK_OK is
; present there, failing the build otherwise.
;
; Must be launched via exec-form RUN (["/usr/bin/supervisord", ...]) so supervisord
; is PID 1 of the build container and actually reproduces the reaping behavior.

[supervisord]
nodaemon=true
user=root
logfile=/dev/stdout
logfile_maxbytes=0
pidfile=/tmp/supervisord-validate.pid
loglevel=info

[program:launch-check]
command=node launch-check.js
directory=/app/src/integrations/whatsapp/bridge
autostart=true
autorestart=false
startsecs=0
startretries=0
exitcodes=0
redirect_stderr=true
stdout_logfile=/tmp/launch-check.log
stdout_logfile_maxbytes=0
environment=CHROMIUM_EXECUTABLE_PATH="/usr/bin/chromium",LAUNCH_CHECK_STOP_PID1="1"
