# supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisord]
nodaemon=true
user=assistant
logfile=/app/logs/supervisord.log
pidfile=/tmp/supervisord.pid
loglevel=info

[program:python-api]
command=/app/.venv/bin/python -m src.main
directory=/app
autostart=true
autorestart=true
startsecs=5
startretries=5
priority=100
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=PATH="/app/.venv/bin:%(ENV_PATH)s",PYTHONPATH="/app"

[program:whatsapp-bridge]
command=/bin/sh /app/src/integrations/whatsapp/bridge/start-bridge.sh
directory=/app/src/integrations/whatsapp/bridge
autostart=true
autorestart=unexpected
exitcodes=0
startsecs=0
startretries=10
priority=200
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=WHATSAPP_BRIDGE_PORT="3001",WHATSAPP_SESSION_DIR="/app/data/whatsapp_session",WHATSAPP_WEBHOOK_URL="http://localhost:8080/api/whatsapp/webhook/incoming",CHROMIUM_EXECUTABLE_PATH="/usr/bin/chromium"
