# supervisord.validate.conf
; Build-time validation config for the WhatsApp bridge Chromium launch.
;
; This is NOT used at runtime. It is run in CI (ci.yml, whatsapp-bridge-check job)
; under `unshare --pid --fork --mount-proc` so supervisord becomes PID 1 and
; reproduces the exact runtime condition: PID 1 reaping Chromium's child processes.
; That is what caused "Failed to launch the browser process: Code: null" when
; --no-zygote was present — a plain `node launch-check.js` invocation would NOT
; catch this regression.
;
; The BRIDGE_DIR environment variable must point to the whatsapp/bridge directory.
; CHROMIUM_EXECUTABLE_PATH must point to the chromium binary to test.
; launch-check.js prints LAUNCH_CHECK_OK on success and signals PID 1 (this
; supervisord) to shut down, so the CI step returns promptly.

[supervisord]
nodaemon=true
user=root
logfile=/dev/stdout
logfile_maxbytes=0
pidfile=/tmp/supervisord-validate.pid
loglevel=info

[program:launch-check]
command=node launch-check.js
directory=%(ENV_BRIDGE_DIR)s
autostart=true
autorestart=false
startsecs=0
startretries=0
exitcodes=0
redirect_stderr=true
stdout_logfile=/tmp/launch-check.log
stdout_logfile_maxbytes=0
environment=LAUNCH_CHECK_STOP_PID1="1"
