# supervisord.validate.conf
; CI-only validation config — NOT used at runtime.
; Runs under `unshare --pid` so supervisord is PID 1, reproducing the container
; condition that triggered "Failed to launch the browser process: Code: null".
; Requires BRIDGE_DIR and CHROMIUM_EXECUTABLE_PATH to be set by the caller.

[supervisord]
nodaemon=true
user=root
logfile=/dev/stdout
logfile_maxbytes=0
pidfile=/tmp/supervisord-validate.pid
loglevel=info

[program:launch-check]
command=node launch-check.js
directory=%(ENV_BRIDGE_DIR)s
autostart=true
autorestart=false
startsecs=0
startretries=0
exitcodes=0
redirect_stderr=true
stdout_logfile=/tmp/launch-check.log
stdout_logfile_maxbytes=0
environment=LAUNCH_CHECK_STOP_PID1="1"
